# lib/dpm-collector.service
[Unit]
Description=DPM-DA510/530 Modbus RTU Collector
After=network.target influxdb.service
Wants=influxdb.service

[Service]
Type=simple
User=dpm
Group=dpm
WorkingDirectory=/opt/dpm-collector
EnvironmentFile=/opt/dpm-collector/.env
ExecStart=/usr/bin/node /opt/dpm-collector/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dpm-collector

[Install]
WantedBy=multi-user.target
